# deploy/supervisor-arl-worker.conf
; ============================================================================
; Supervisor - Worker de cola para Gestión de Afiliaciones ARL
; ============================================================================
;
; Procesa los jobs en segundo plano (ej. aprobación masiva de certificados ARL
; en lotes grandes con App\Jobs\ProcesarCertificadosArlJob).
;
; INSTALACIÓN (servidor Linux):
;   1. sudo apt install supervisor
;   2. Copiar este archivo a: /etc/supervisor/conf.d/arl-worker.conf
;   3. Ajustar las rutas, el usuario y la cantidad de procesos según el servidor.
;   4. sudo supervisorctl reread
;      sudo supervisorctl update
;      sudo supervisorctl start arl-worker:*
;
; REQUISITOS PREVIOS en .env:
;   QUEUE_CONNECTION=database     (o redis)
;   Y haber corrido: php artisan migrate   (crea la tabla jobs)
;
; TRAS CADA DESPLIEGUE:
;   php artisan queue:restart     (para que el worker tome el código nuevo)
; ============================================================================

[program:arl-worker]
process_name=%(program_name)s_%(process_num)02d
; Ajustar la ruta absoluta del proyecto en el servidor:
command=php /var/www/gestion-afiliaciones-arl/artisan queue:work --sleep=3 --tries=1 --timeout=3600 --max-time=3600
directory=/var/www/gestion-afiliaciones-arl
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
; Usuario del servidor web (en muchos servidores es www-data):
user=www-data
; Cantidad de workers en paralelo (1 es suficiente para lotes ARL):
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/gestion-afiliaciones-arl/storage/logs/worker.log
stopwaitsecs=3600
